# ansible/templates/10-poweroff.rules.j2
ACTION=="remove", ENV{ID_VENDOR_ID}=={{ usb_vendor_id }}, ENV{ID_MODEL_ID}=={{ usb_model_id }}, RUN+="/sbin/shutdown -h now"